# torrserver.service
[Unit]
Description=TorrServer

[Service]
Type=simple
ExecStart=/bin/torrserver

[Install]
WantedBy=multi-user.target

